# home/.config/systemd/user/dms.service.d/override.conf
[Service]
Environment=DMS_LOCAL_PATH=%h/.config/quickshell/dms
